# systemd/.config/systemd/user/deckapp.service
[Unit]
Description=Deckapp

[Service]
ExecStart=/home/alice/work/incubator/deck/packages/example/bin/index.js

[Install]
WantedBy=default.target
